# config/system/etc/rc.local
#!/bin/bash

wait-for-docker
#Creating traefik reverse proxy network
echo "Hello there!"
